# imbu/conf/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".

[unix_http_server]
file=%(here)s/../imbu-supervisor.sock   ; (the path to the socket file)

[inet_http_server]
port=127.0.0.1:9001

[supervisord]
environment=CORTICAL_API_KEY="Get Key from http://www.cortical.io/resources_apikey.html"
pidfile=%(here)s/../imbu-supervisord.pid
identifier=imbu-supervisor
logfile=%(here)s/../logs/imbu-supervisord.log
# NOTE: logfile_maxbytes=0 turns off supervisor log rotation to prevent conflict
# with Taurus' higher-level log rotation triggered by crontab
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://:

;*************** FLUENT-API **************
[program:fluent-api]
command=uwsgi --enable-threads --socket 0.0.0.0:19002 --master --vacuum --idle 300 --processes 2 --threads 4 --wsgi-file engine/fluent_api.py
process_name=%(program_name)s_%(process_num)02d
directory=%(here)s/..
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=20
stdout_logfile=logs/uwsgi.log
redirect_stderr=true
stopsignal=INT

;*************** IMBU-GROUP **************
[group:imbu]
programs=fluent-api
