# infrastructure/ami-tools/boot-scripts/centos7/grok-stopped.service
[Unit]
Description=GROK stopped
 
[Install]
Alias=grok-stopped.service
 
[Service]
Type=forking
ExecStart=/etc/init.d/grok-preload grok-stopped
ExecStop=/etc/init.d/grok-preload stop
ExecReload=/etc/init.d/grok-preload reload
Restart=on-abort

