# infrastructure/ami-tools/boot-scripts/centos7/rabbitmq-server.service
[Unit]
Description=RabbitMQ
After=network.target
After=syslog.target

[Install]
Alias=rabbitmq-server.service
WantedBy=multi-user.target

[Service]
ExecStart=/usr/sbin/rabbitmq-server
ExecStop=/usr/sbin/rabbitmqctl stop
Restart=always
