# infrastructure/ami-tools/boot-scripts/centos7/rabbitmq-server.service
[Unit]
Description=RabbitMQ

[Install]
Alias=rabbitmq.service

[Service]
ExecStart=/usr/sbin/rabbitmq-server
ExecStop=/usr/sbin/rabbitmqctl stop
Restart=always
